# packaging/debian_arm/lib/systemd/system/lake.service
[Unit]
Description=openbank lake message relay
After=inetd.service

[Service]
Type=simple
EnvironmentFile=/etc/init/lake.conf
StartLimitBurst=0
ExecStart=/openbank/services/lake/lake
ExecStop=/openbank/services/lake/lake-stop
StandardInput=null
LimitNOFILE=1048576
LimitNPROC=infinity
LimitCORE=infinity
Restart=on-failure
RestartSec=2
TimeoutStartSec=0
KillMode=process

[Install]
WantedBy=multi-user.target
