# packaging/debian_amd64/lib/systemd/system/lake.path
[Unit]
Description=openbank lake services namespace
PartOf=lake.service

[Path]
PathChanged=/etc/init/lake.conf

[Install]
WantedBy=lake.service
