# packaging/debian_amd64/lib/systemd/system/lake.path
[Unit]
Description=openbank lake services namespace
PartOf=lake.service
ConditionPathExists=/etc/init/lake.conf

[Path]
PathChanged=/etc/init/lake.conf
PathModified=/etc/init/lake.conf
Unit=lake.service

[Install]
WantedBy=lake.service
