# packaging/debian/lake-watcher.path
[Unit]
Description=openbank lake configuration watcher
After=lake.service

[Path]
PathChanged=/etc/lake/conf.d
Unit=lake-watcher.service

[Install]
WantedBy=paths.target
