# packaging/debian/lake.path
[Unit]
Description=openbank lake services namespace
PartOf=lake.service
ConditionPathExists=/etc/lake/conf.d/init.conf

[Path]
PathChanged=/etc/lake/conf.d/init.conf
PathModified=/etc/lake/conf.d/init.conf
Unit=lake.service

[Install]
WantedBy=lake.service
